# Project/linux_script/setup/11-local.rules
ACTION=="add", SUBSYSTEMS=="usb", DRIVERS=="*hid*", RUN+="/home/vince/bin/usb_add.sh"
